# teveclub.service
[Unit]
Description=Teveclub Django Application
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/teveclub/django
Environment="PATH=/home/ubuntu/teveclub/venv/bin"
Environment="PYTHONPATH=/home/ubuntu/teveclub/django"
ExecStart=/home/ubuntu/teveclub/venv/bin/gunicorn teveclub_web.wsgi:application \
    --bind 0.0.0.0:8000 \
    --workers 3 \
    --timeout 60 \
    --access-logfile /home/ubuntu/teveclub/logs/access.log \
    --error-logfile /home/ubuntu/teveclub/logs/error.log \
    --log-level info
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
